# rootfs/etc/s6-overlay/s6-rc.d/testing/run
#!/command/with-contenv bash

nginx -g "daemon off;"
