# docker/rootfs/etc/s6-overlay/s6-rc.d/api/run
#!/command/execlineb -P
cd /app/api
./scripts/start.sh